# templates/homarr.service.j2
# =============================================================================
# Homarr Dashboard Systemd Service
# =============================================================================
#
# This systemd service file manages the Homarr dashboard application.
# Homarr is a customizable home page / dashboard with service integrations.
#
# Service Details:
#   - Runs as {{ homarr_user }} user
#   - Data stored in {{ homarr_data_dir }}
#   - Web interface on port {{ homarr_web_port }}
#
# Author: Generated by Ansible
# =============================================================================

[Unit]
Description=Homarr Dashboard
After=network.target

[Service]
Type=simple
User={{ homarr_user }}
Group={{ homarr_group }}
WorkingDirectory={{ homarr_install_dir }}
ExecStart=/usr/bin/node {{ homarr_install_dir }}/node_modules/next/dist/bin/next start -p {{ homarr_web_port }}
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=homarr

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ homarr_data_dir }} {{ homarr_install_dir }}

# Environment
Environment="NODE_ENV=production"
Environment="PORT={{ homarr_web_port }}"
EnvironmentFile={{ homarr_config_file }}

[Install]
WantedBy=multi-user.target
