# templates/lets-encrypt-renew.timer.j2
# =============================================================================
# Let's Encrypt Renewal Systemd Timer Unit File
# =============================================================================
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY! ⚠️
# Any manual changes will be overwritten on the next Ansible run.
# Edit the template file or Ansible variables instead.
# -----------------------------------------------------------------------------
# File: roles/networking/templates/lets-encrypt-renew.timer.j2
# Deployed: /etc/systemd/system/lets-encrypt-renew.timer
# Generated by Ansible on: {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
# -----------------------------------------------------------------------------
#
# Description:
#   Systemd timer unit to schedule Let's Encrypt certificate renewal using Certbot.
#   The renewal time is randomized per host between 01:00 and 04:59 to avoid thundering herd issues.
#   Triggers the lets-encrypt-renew.service unit on a monthly schedule.
#
# Usage:
#   sudo systemctl enable lets-encrypt-renew.timer    # Enable timer on boot
#   sudo systemctl start lets-encrypt-renew.timer     # Start timer
#   sudo systemctl status lets-encrypt-renew.timer    # Check timer status
#   sudo systemctl list-timers                        # List all active timers
#
# Dependencies:
#   - Requires lets-encrypt-renew.service
#
# References:
#   - https://www.freedesktop.org/software/systemd/man/systemd.timer.html
#   - https://certbot.eff.org/
#   - https://letsencrypt.org/docs/
# =============================================================================

[Unit]
Description="Timer for Let's Encrypt certificate renewal"

[Timer]
OnCalendar=*-*-01 {{ lets_encrypt_random_renew_time }}
Persistent=true

[Install]
WantedBy=timers.target