# templates/lets-encrypt-renew.service.j2
# =============================================================================
# Let's Encrypt Renewal Systemd Service Unit File
# =============================================================================
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY! ⚠️
# Any manual changes will be overwritten on the next Ansible run.
# Edit the template file or Ansible variables instead.
# -----------------------------------------------------------------------------
# File: roles/networking/templates/lets-encrypt-renew.service.j2
# Deployed: /etc/systemd/system/lets-encrypt-renew.service
# Generated by Ansible on: {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
# -----------------------------------------------------------------------------
#
# Description:
#   Systemd unit file to automate Let's Encrypt certificate renewal using Certbot.
#   Includes a random sleep to avoid thundering herd issues with Let's Encrypt servers.
#
# Usage:
#   sudo systemctl enable lets-encrypt-renew.service    # Enable on boot
#   sudo systemctl start lets-encrypt-renew.service     # Start service
#   sudo systemctl status lets-encrypt-renew.service    # Check status
#   sudo journalctl -u lets-encrypt-renew.service       # View logs
#   sudo journalctl -u lets-encrypt-renew.service -f    # View and follow logs
#
# Security:
#   - Runs as limited-privilege user: {{ lets_encrypt_user }}
#   - Group: {{ lets_encrypt_group }}
#
# Dependencies:
#   - Requires systemd timer (see lets-encrypt-renew.timer)
#
# References:
#   - https://certbot.eff.org/
#   - https://letsencrypt.org/docs/
# =============================================================================

[Unit]
Description="Renew Let's Encrypt certificates with random sleep"
Documentation="https://certbot.eff.org/instructions?ws=other&os=pip&tab=wildcard"

After=network-online.target

[Service]
Type=oneshot
User={{ lets_encrypt_user }}
Group={{ lets_encrypt_group }}

# Renew the certificates quietly
# -q: quiet mode, only output errors
ExecStart={{ lets_encrypt_certbot_binary }} renew  -q \
                                        --config-dir {{ lets_encrypt_config_dir }} \
                                        --work-dir {{ lets_encrypt_config_dir }} \
                                        --logs-dir {{ lets_encrypt_config_dir }}/logs

# Log to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lets_encrypt_renew