# templates/nfs-systemd.mount.j2
# ==========================================================================
# Systemd Mount Unit for NFS Share
# ==========================================================================
#
# This unit file mounts an NFS share using systemd.
# Generated by Ansible for {{ item.server }}:{{ item.share_path }}
#
# https://www.freedesktop.org/software/systemd/man/latest/systemd.mount.html
#
# Author: Generated by Ansible
# ==========================================================================

[Unit]
Description=Mount NFS share {{ item.share_path }} from {{ item.server }}
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/etc/nsswitch.conf

[Mount]
What={{ item.server }}:{{ item.share_path }}
Where={{ item.mount_path }}
Type=nfs
{# Options={{ item.options | default(nfs_default_options) }},_netdev,x-systemd.automount,x-systemd.device-timeout=15s,x-systemd.mount-timeout=30s #}
Options=vers=3,soft,timeo=60,rsize=32768,wsize=32768,nolock,_netdev,x-systemd.automount,x-systemd.device-timeout=60s,x-systemd.mount-timeout=60s

TimeoutIdleSec=300
TimeoutSec=30

[Install]
WantedBy=multi-user.target