# templates/openvpn-client-systemd.service.j2
# =============================================================================
# OpenVPN Client Systemd Service Unit File
# =============================================================================
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY! ⚠️
# Any manual changes will be overwritten on the next Ansible run.
# Edit the template file or Ansible variables instead.
# -----------------------------------------------------------------------------
# File: roles/networking/templates/openvpn-client-systemd.service.j2
# Deployed: /etc/systemd/system/openvpn@client.service (or openvpn@<name>.service)
# Generated by Ansible on: {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
# -----------------------------------------------------------------------------
#
# Description:
#   Systemd unit file to manage OpenVPN client tunnels.
#   Supports dynamic configuration via instance name (e.g., openvpn@client).
#   Runs a pre-start script to update dynamic VPN server rules if configured.
#
# Usage:
#   sudo systemctl enable openvpn@client.service    # Enable on boot
#   sudo systemctl start openvpn@client.service     # Start service
#   sudo systemctl status openvpn@client.service    # Check status
#   sudo journalctl -u openvpn@client.service       # View logs
#   sudo journalctl -u openvpn@client.service -f    # View and follow logs
#
# References:
#   - https://community.openvpn.net/openvpn/wiki/Systemd
#   - https://man7.org/linux/man-pages/man5/systemd.service.5.html
# =============================================================================

[Unit]
Description=OpenVPN tunnel for %I
After=network-online.target
Wants=network-online.target
Documentation=man:openvpn(8)
Documentation=https://community.openvpn.net/openvpn/wiki/Openvpn24ManPage
Documentation=https://community.openvpn.net/openvpn/wiki/HOWTO

[Service]
Type=notify
PrivateTmp=true
WorkingDirectory=/etc/openvpn/client
ExecStartPre={{ vpn_server_dynamic_update_script }}
ExecStart=/usr/sbin/openvpn --suppress-timestamps --nobind --config %i.conf
CapabilityBoundingSet=CAP_IPC_LOCK CAP_NET_ADMIN CAP_NET_RAW CAP_SETGID CAP_SETUID CAP_SETPCAP CAP_SYS_CHROOT CAP_DAC_OVERRIDE
LimitNPROC=10
DeviceAllow=/dev/null rw
DeviceAllow=/dev/net/tun rw
ProtectSystem=true
ProtectHome=true
KillMode=process

[Install]
WantedBy=multi-user.target