# templates/ssh-systemd.service.j2
# =============================================================================
# OpenSSH Systemd Service File - Ansible Managed
# =============================================================================
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY!
# Any manual changes will be overwritten on the next playbook run.
# Edit the template file or Ansible variables instead.
#
# Template path: roles/networking/templates/ssh-systemd.service.j2
# Deploy path: {{ ssh_service_file }}
# Last deployed: {{ ansible_date_time.iso8601 }}
# =============================================================================
#
# This systemd service file configures OpenSSH to run as a system service,
# providing secure remote access.
#
# Features:
#   - Waits for network to be online before starting
#   - Uses auditd for logging
#   - Hardened with systemd security options
#
# Usage:
#   sudo systemctl enable ssh
#   sudo systemctl start ssh
#   sudo systemctl status ssh
#   sudo journalctl -u ssh -f
#
# Security Notes:
#   - ssh runs with limited privileges
#
# Dependencies:
#   - Requires network-online.target and auditd.service
#
# References:
#   - man:sshd(8)
#   - man:sshd_config(5)
#
# Generated by Ansible on {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
# ---

[Unit]
Description=OpenBSD Secure Shell server
Documentation=man:sshd(8) man:sshd_config(5)
#After=network.target auditd.service

# Service Dependencies
######################

# This service requires systemd-networkd to be online before starting as we are 
# biding to the network interface ip address for security reasons.
# NOTE: A networkd inteface profile is required for @interface.service is being used
#After=network-online.target systemd-networkd-wait-online@{{ ssh_primary_network_interface }}.service auditd.#service
#Wants=systemd-networkd-wait-online@{{ ssh_primary_network_interface }}.service
After=network-online.target auditd.service

# Condition to prevent starting if sshd is not to be run
# This allows for disabling SSH without removing the service file
ConditionPathExists=!/etc/ssh/sshd_not_to_be_run


[Service]
EnvironmentFile=-/etc/default/ssh
ExecStartPre=/usr/sbin/sshd -t
ExecStart=/usr/sbin/sshd -D $SSHD_OPTS
ExecReload=/usr/sbin/sshd -t
ExecReload=/bin/kill -HUP $MAINPID

Type=notify
KillMode=process
Restart=on-failure
RestartPreventExitStatus=255

# Directory creation and permissions
RuntimeDirectory=sshd
RuntimeDirectoryMode=0755

# Logging to systemd
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ssh


[Install]
WantedBy=multi-user.target
Alias=sshd.service