# templates/vpn-server-update.service.j2
# =============================================================================
# VPN Server Dynamic Update Systemd Service Unit File
# =============================================================================
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY! ⚠️
# Any manual changes will be overwritten on the next Ansible run.
# Edit the template file or Ansible variables instead.
# -----------------------------------------------------------------------------
# File: roles/networking/templates/vpn-server-update.service.j2
# Deployed: /etc/systemd/system/vpn-server-update.service
# Generated by Ansible on: {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
# -----------------------------------------------------------------------------
#
# Description:
#   Systemd unit file to automate dynamic VPN server rule updates in /etc/ufw/before.rules.
#   Runs a custom bash script to update firewall rules for VPN kill switch.
#
# Usage:
#   sudo systemctl enable vpn-server-update.service    # Enable on boot
#   sudo systemctl start vpn-server-update.service     # Start service
#   sudo systemctl status vpn-server-update.service    # Check status
#   sudo journalctl -u vpn-server-update.service       # View logs
#   sudo journalctl -u vpn-server-update.service -f    # View and follow logs
#
# Dependencies:
#   - Requires systemd timer (see vpn-server-update.timer)
#   - Requires vpn-server-update.sh script deployed and executable
#
# References:
#   - https://wiki.archlinux.org/title/Systemd/Timers
#   - https://man7.org/linux/man-pages/man5/systemd.service.5.html
# =============================================================================

[Unit]
Description="Custom bash script to update the dynamic VPN server configuration in /etc/ufw/before.rules"
After=network-online.target
PartOf=openvpn@client.service
SyslogIdentifier=ovpn-client

[Service]
Type=oneshot
User=root
Group=root

# Update the VPN server configuration
ExecStart={{ vpn_server_dynamic_update_script }}