# templates/vpn-server-update.timer.j2
# =============================================================================
# VPN Server Dynamic Update Systemd Timer Unit File
# =============================================================================
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY! ⚠️
# Any manual changes will be overwritten on the next Ansible run.
# Edit the template file or Ansible variables instead.
# -----------------------------------------------------------------------------
# File: roles/networking/templates/vpn-server-update.timer.j2
# Deployed: /etc/systemd/system/vpn-server-update.timer
# Generated by Ansible on: {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
# -----------------------------------------------------------------------------
#
# Description:
#   Systemd timer unit to schedule dynamic VPN server rule updates every 15 minutes.
#   Triggers the vpn-server-update.service unit on a recurring schedule.
#
# Usage:
#   sudo systemctl enable vpn-server-update.timer    # Enable timer on boot
#   sudo systemctl start vpn-server-update.timer     # Start timer
#   sudo systemctl status vpn-server-update.timer    # Check timer status
#   sudo systemctl list-timers                      # List all active timers
#
# Dependencies:
#   - Requires vpn-server-update.service
#
# References:
#   - https://www.freedesktop.org/software/systemd/man/systemd.timer.html
# =============================================================================

[Unit]
Description="Timer for VPN Server Dynamic Update"

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target